# microservice-base/supervisor/supervisor.d/dns-updater-first.conf
[program:dns-updater-first]
priority = 45
command=/usr/bin/resolve-dns
user = root
autostart = true
autorestart = false
