# microservice-base/supervisor/supervisor.d/dns-updater-first.conf
[program:dns-updater-first]
priority = 50
command=/usr/bin/resolve-dns
user = root
umask = 022
autostart = true
autorestart = unexpected
exitcodes = 0,2
startsecs = 0
