# base-image-centos/supervisor/supervisor.d/entrypoint.conf
[program:entrypoint]
priority = 30
command=run-parts /etc/entrypoint.d
user = root
autorestart = false
