# base-consul/supervisor/supervisor.d/consul.conf
[program:consul]
priority = 50

user = root
command = run-consul.sh agent %(ENV_CONSUL_PARAMS)s
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
redirect_stderr = true
stopwaitsecs = 20
