# consul-cassandra/supervisor/supervisor.d/cassandra.conf
[program:cassandra]
priority = 60

user = root
command = /run-cassandra.sh -f
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
redirect_stderr = true
stopwaitsecs = 600
