# base-image/supervisor/supervisor.d/sshd.conf
[program:sshd]
priority = 20
command = /usr/sbin/sshd -D
autostart = true
autorestart = true
user = root
