# consul-nomad/supervisor/supervisor.d/consul.conf
[program:consul]
priority = 40

user = root
command = /run-consul.sh agent
#command = bash -c '/bin/consul agent -config-dir=/consul/config -data-dir=/consul/data'
