# consul-nomad/supervisor/supervisor.d/consul.conf
[program:consul]
priority = 40

user = nobody
command = bash -c '/bin/consul agent -config-dir=/consul/config -data-dir=/consul/data'
