# consul-nomad/supervisor/supervisor.d/nomad.conf
[program:nomad]
priority = 45

user = nobody
command = bash -c '/bin/nomad agent -config=/nomad/config -data-dir=/nomad-data'
