# consul-mariadb/supervisor/supervisor.d/mariadb.conf
[program:mariadb]
priority = 60

user = root
command = /run-mariadb.sh
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
redirect_stderr = true
