# consul-mysql/supervisor/supervisor.d/mysql.conf
[program:mysql]
priority = 60

user = root
command = run-mysql.sh
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
redirect_stderr = true
stopwaitsecs = 60
